# thin-client-assistant/systemd/pve-thin-client-prepare.service
[Unit]
Description=Prepare PVE thin-client runtime
After=local-fs.target
Before=getty@tty1.service

[Service]
Type=oneshot
ExecStart=/usr/local/lib/pve-thin-client/prepare-runtime.sh

[Install]
WantedBy=multi-user.target
