# thin-client-assistant/live-build/config/includes.chroot/etc/profile.d/pve-thin-client-autostart.sh
# Startup is handled by dedicated systemd services bound to tty1.
return 0 2>/dev/null || exit 0
