# proxmox-host/systemd/beagle-ui-reapply.service
[Unit]
Description=Reapply Beagle OS Proxmox UI hook
After=local-fs.target pveproxy.service
Wants=pveproxy.service
ConditionPathExists=/usr/share/pve-manager/index.html.tpl

[Service]
Type=oneshot
EnvironmentFile=-/etc/beagle/host.env
ExecStart=__INSTALL_DIR__/scripts/install-proxmox-ui-integration.sh
WorkingDirectory=__INSTALL_DIR__

[Install]
WantedBy=multi-user.target
