# proxmox-host/systemd/beagle-ui-reapply.path
[Unit]
Description=Watch Proxmox UI files and reapply Beagle OS after updates

[Path]
PathChanged=/usr/share/pve-manager/index.html.tpl
PathChanged=/usr/share/pve-manager/js/pvemanagerlib.js
Unit=beagle-ui-reapply.service

[Install]
WantedBy=multi-user.target
