# proxmox-host/systemd/beagle-control-plane.service
[Unit]
Description=Beagle OS Control Plane API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/beagle/beagle-manager.env
ExecStart=/usr/bin/python3 __INSTALL_DIR__/proxmox-host/bin/beagle-control-plane.py
Restart=always
RestartSec=3
User=root
Group=root

[Install]
WantedBy=multi-user.target
