# proxmox-host/systemd/beagle-artifacts-refresh.service
[Unit]
Description=Refresh Beagle OS hosted artifacts
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/beagle/host.env
ExecStart=__INSTALL_DIR__/scripts/refresh-host-artifacts.sh
WorkingDirectory=__INSTALL_DIR__

[Install]
WantedBy=multi-user.target
