# proxmox-host/systemd/pve-dcv-artifacts-refresh.timer
[Unit]
Description=Periodic refresh of PVE DCV Integration hosted artifacts

[Timer]
OnBootSec=10min
OnUnitActiveSec=7d
AccuracySec=15min
Persistent=true
Unit=pve-dcv-artifacts-refresh.service

[Install]
WantedBy=timers.target
